# init.scripts/systemd/snapclient.service
[Unit]
Description=Snapclient
After=network.target

[Service]
ExecStart=/usr/sbin/snapclient

[Install]
WantedBy=multi-user.target 
