# init.scripts/systemd/snapserver.service
[Unit]
Description=Snapserver
After=network.target

[Service]
ExecStart=/usr/sbin/snapserver

[Install]
WantedBy=multi-user.target 
